# templates/systemd/minecraft.service.j2
{#
SPDX-FileCopyrightText: 2023 Slavi Pantaleev
SPDX-FileCopyrightText: 2024 Nikita Chernyi
SPDX-FileCopyrightText: 2025 XHawk87
SPDX-FileCopyrightText: 2025, 2026 Suguru Hirahara

SPDX-License-Identifier: AGPL-3.0-or-later
#}

[Unit]
Description=Minecraft ({{ minecraft_identifier }})
{% for service in minecraft_systemd_required_services_list %}
Requires={{ service }}
After={{ service }}
{% endfor %}
{% for service in minecraft_systemd_wanted_services_list %}
Wants={{ service }}
{% endfor %}
DefaultDependencies=no

[Service]
Type=simple
Environment="HOME={{ devture_systemd_docker_base_systemd_unit_home_path }}"
ExecStartPre={{ devture_systemd_docker_base_host_command_sh }} -c '{{ minecraft_container_cleanup_command | trim | replace("\n", " ") }}'

ExecStartPre={{ devture_systemd_docker_base_host_command_docker }} create \
    --rm \
    --name={{ minecraft_identifier }} \
    --log-driver=none \
    --user={{ minecraft_uid }}:{{ minecraft_gid }} \
    --cap-drop=ALL \
    --network={{ minecraft_container_network }} \
    {% if minecraft_container_tcp_host_bind_port %}
    -p {{ minecraft_container_tcp_host_bind_port }}:{{ minecraft_container_tcp_port }} \
    {% endif %}
    {% for cap in minecraft_security_allow_capabilities %}
    --cap-add={{ cap }} \
    {% endfor %}
    {% if not minecraft_security_write_access %}
    --read-only \
    {% endif %}
    {% if minecraft_security_privileged %}
    --privileged \
    {% endif %}
    {% for dns in minecraft_container_dns_servers %}
    --dns={{ dns }} \
    {% endfor %}
    --env-file={{ minecraft_base_path }}/env \
    --label-file={{ minecraft_base_path }}/labels \
    --health-interval={{ minecraft_container_health_interval }} \
    -v {{ minecraft_data_path }}:/data \
    {% for volume in minecraft_container_additional_volumes %}
    --mount type={{ volume.type | default('bind' if '/' in volume.src else 'volume') }},src={{ volume.src }},dst={{ volume.dst }}{{ (',' + volume.options) if volume.options | default('') else '' }} \
    {% endfor %}
    {% for arg in minecraft_container_extra_arguments %}
    {{ arg }} \
    {% endfor %}
    {{ minecraft_container_image }}

{% for network in minecraft_container_additional_networks %}
ExecStartPre={{ devture_systemd_docker_base_host_command_docker }} network connect {{ network }} {{ minecraft_identifier }}
{% endfor %}

ExecStart={{ devture_systemd_docker_base_host_command_docker }} start --attach {{ minecraft_identifier }}

ExecStop=-{{ devture_systemd_docker_base_host_command_sh }} -c '{{ devture_systemd_docker_base_host_command_docker }} stop -t {{ devture_systemd_docker_base_container_stop_grace_time_seconds }} {{ minecraft_identifier }} 2>/dev/null || true'
ExecStop=-{{ devture_systemd_docker_base_host_command_sh }} -c '{{ devture_systemd_docker_base_host_command_docker }} rm {{ minecraft_identifier }} 2>/dev/null || true'

Restart=always
RestartSec=30
SyslogIdentifier={{ minecraft_identifier }}

[Install]
WantedBy=multi-user.target
